# system_monitor.service
[Unit]
Description=Linux System Monitor Publisher
Documentation=man:system_monitor(1)

[Service]
Type=simple
ExecStart=/usr/bin/system_monitor -c /etc/sensors/system_monitor.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target

